# scripts/corby.service
[Service]
WorkingDirectory=/root/
ExecStart=/root/corby/scripts/start.sh
StartLimitBurst=0